# weather-bot.service
[Unit]
Description=Polymarket Weather Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/PollyMorph
ExecStart=/usr/bin/python3 -m tools.weather.executor --loop --live --probe --bankroll 164
Restart=always
RestartSec=60
StandardOutput=append:/root/PollyMorph/logs/weather-bot.log
StandardError=append:/root/PollyMorph/logs/weather-bot.log

# Give the bot time to finish a scan before killing it on stop
TimeoutStopSec=120
KillSignal=SIGTERM

# Environment — load from .env file in the project root
EnvironmentFile=/root/PollyMorph/.env
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
